# supervisord.conf
[supervisord]
nodaemon=true

[program: bundle]
command=bundle install && sleep 20
autorestart=false

[program:dashing]
command=service dashboard start
autorestart=false
