# supervisord.conf
[supervisord]
nodaemon=true

[program:dashing]
command=service dashboard start
autorestart=false
